# scripts/launchd/grace.change-executor-batch.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.grace.change-executor-batch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/python@3.14/bin/python3.14</string>
        <string>-m</string>
        <string>src.ontology.change_executor</string>
        <string>batch</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/grace</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/YOUR_USERNAME/grace</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/grace-change-executor-batch.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/grace-change-executor-batch.stderr.log</string>
</dict>
</plist>
