# scripts/launchd/grace.confidence-decay.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Example launchd plist for the GrACE confidence decay batch (Chunk 35a, D264).

  Operational model (D246 mirror): the decay batch is CLI-only. APScheduler
  / FastAPI lifespan integration is intentionally omitted. Operators
  schedule the CLI via host cron / launchd.

  Install (per-user):
    1. Copy this file to ~/Library/LaunchAgents/com.grace.confidence-decay.plist
    2. Edit the placeholders:
       - GRACE_HOME       -> REPLACE THIS: absolute path to your grace checkout
       - PYTHON_BIN       -> REPLACE THIS: absolute path to the Python you used
                             to install grace dependencies (run `python3 -m pip -V`
                             and use the bin directory it reports)
    3. Load it:           launchctl load ~/Library/LaunchAgents/com.grace.confidence-decay.plist
    4. Verify:            launchctl list | grep com.grace.confidence-decay

  Schedule below runs once weekly on Sundays at 03:00 local time. The
  CLI exits 0 on success, 2 on argparse / config errors, non-zero on
  ArcadeDB failure. launchd records exits in the system log.

  Scope (v3 B1 resolution): the batch reads from and writes to ArcadeDB
  ONLY. It does not touch Postgres `extraction_claims`.

  The observation-time argument is filled at run time so each
  invocation is idempotent under a fixed clock anchor (D264
  idempotency contract).
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.grace.confidence-decay</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/PYTHON_BIN -m src.extraction.confidence_decay --observation-time "$(date -u +%Y-%m-%dT%H:%M:%SZ)" --config config/decay_config.yaml</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/GRACE_HOME</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/GRACE_HOME</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/GRACE_HOME/.build-state/grace-confidence-decay.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/GRACE_HOME/.build-state/grace-confidence-decay.stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
