# scripts/launchd/grace.agent-daemon.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.grace.agent-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/python@3.14/bin/python3.14</string>
        <string>-m</string>
        <string>src.ontology.agent_daemon</string>
        <string>run</string>
        <string>--tick-interval</string>
        <string>300</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/grace</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/YOUR_USERNAME/grace</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/grace-agent-daemon-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/grace-agent-daemon-stderr.log</string>
</dict>
</plist>
