# scripts/launchd/grace.voice-tone-profiling.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.grace.voice-tone-profiling</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/python@3.14/bin/python3.14</string>
        <string>-m</string>
        <string>src.ingestion.communications.voice_tone</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/grace</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/grace-voice-tone-profiling.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/grace-voice-tone-profiling.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/YOUR_USERNAME/grace</string>
    </dict>
</dict>
</plist>
